# deploy/whitemagic-dashboard.service
[Unit]
Description=WhiteMagic PWA Dashboard
After=network.target whitemagic-api.service

[Service]
Type=simple
User=whitemagic
WorkingDirectory=/home/whitemagic/WHITEMAGIC/apps/site

# Node.js environment
Environment="NODE_ENV=production"
Environment="PORT=3002"
Environment="PATH=/usr/bin:/usr/local/bin"

ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/whitemagic/WHITEMAGIC/apps/site/.next

[Install]
WantedBy=multi-user.target
